# modules/defaults/files/etc/profile.d/prompt.sh
#
# This file is maintained by Puppet.
# Changes will be overwritten.
#
export PS1='\u@\h:\w\$'
