# modules/defaults/files/etc/systemd/system/drtm-daily.service
#
# Date: 2017-06-04
#
[Unit]
Description=Dartmore daily

[Service]
Type=oneshot
ExecStart=/opt/drtm/sbin/drtm daily
