# modules/defaults/files/etc/systemd/system/drtm-daily.service
[Unit]
Description=Dartmore daily

[Service]
Type=oneshot
ExecStart=/opt/drtm/sbin/drtm daily
